# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:hermes]
command=/app/start.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:syncthing]
# -home points to our persistent volume path
# -gui-address ensures we can access the UI internally if needed
command=/usr/bin/syncthing -home=/data/syncthing -gui-address=0.0.0.0:8384 -no-browser
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
