# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:hermes]
command=/app/hermes_start.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:syncthing-config]
command=/app/syncthing_config.sh
priority=1
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:syncthing]
# -home points to our persistent volume path
# -gui-address ensures we can access the UI internally if needed
command=/usr/bin/syncthing -home=/data/syncthing -gui-address=0.0.0.0:8384 -no-browser
priority=2
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gbrain-init]
# Initialize GBrain with Postgres database and verify setup
command=/app/gbrain_start.sh
priority=3
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gbrain-sync]
# Live sync: sync brain repo and embed stale documents (every 15 minutes)
command=/bin/bash -c "while true; do /app/gbrain_sync.sh; sleep 900; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gbrain-dream]
# Dream cycle: nightly memory consolidation and pattern detection (runs at 2 AM UTC)
command=/bin/bash -c "while true; do sleep $(( ($(date -d '02:00 UTC' +%s) - $(date +%s) + 86400) % 86400 )); /app/gbrain_dream.sh; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
